# package/openwsman.service.in
[Unit]
Description=Openwsman WS-Management Service
After=syslog.target

[Service]
Type=forking
ExecStart=@prefix@/sbin/openwsmand -S

[Install]
WantedBy=multi-user.target
